# crontab
30 22 * * * test -e /root/rsync_backup/cron_mirror_rsync.sh && /root/rsync_backup/cron_mirror_rsync.sh

